# autopatch/feishu-autopatch.timer
[Unit]
Description=飞书 asar 自动重打(开机+定时兜底)

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min
Unit=feishu-autopatch.service

[Install]
WantedBy=timers.target
